# deployment_scripts/puppet/patches/zookeeper/templates/zookeeper.service.erb
# Autogenerated by Puppet from zookeeper/templates/zookeeper.service.erb
# Do not edit

[Unit]
Description=Apache ZooKeeper
<%if @systemd_unit_want -%>
Want=<%=@systemd_unit_want %>
<% end -%>
<%if @systemd_unit_after -%>
After=<%=@systemd_unit_after %>
<% end -%>
ConditionPathExists=/var/lib/zookeeper/myid

[Service]
Type=simple
User=<%= @user %>
Group=<%= @group %>
SyslogIdentifier=zookeeper
EnvironmentFile=/etc/zookeeper/conf/environment
EnvironmentFile=-/etc/default/zookeeper
ExecStart=/usr/bin/java -cp ${CLASSPATH} $JAVA_OPTS -Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.local.only=${JMXLOCALONLY} -Dzookeeper.log.dir=${ZOO_LOG_DIR} -Dzookeeper.root.logger=${ZOO_LOG4J_PROP} $ZOOMAIN $ZOOCFG
SuccessExitStatus=143
# ZooKeeper is "fail-fast", see https://zookeeper.apache.org/doc/r3.4.8/zookeeperAdmin.html#sc_supervision
Restart=always

[Install]
WantedBy=multi-user.target